# led_mon_server/led-monitor.service
[Unit]
Description=LED-Monitor (Server)
After=syslog.target

[Service]
Type=simple
WorkingDirectory=/etc/led_mon/
ExecStart=/etc/led_mon/server.py
SyslogIdentifier=led-monitor
StandardOutput=syslog
StandardError=syslog

[Install]
WantedBy=multi-user.target

